# dockerfiles/js/root/etc/services.d/node/run
#!/usr/bin/with-contenv bash

if [[ 'development' == $NODE_ENV ]]
then
    exec npm run develop --prefix /app
else
    exec npm start --prefix /app
fi
